# shopixy-backup.timer
[Unit]
Description=Run Shopixy Backup Daily

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
